# conf/supervisor/tno.ini
[program:skydns]
command=/webapp/conf/supervisor/skydns.sh
directory=/webapp/conf/supervisor/
autostart=true
autorestart=true
stopsignal=TERM
priority=100

[program:skydock]
command=/webapp/conf/supervisor/skydock.sh
directory=/webapp/conf/supervisor/
autostart=true
autorestart=true
stopsignal=TERM
priority=200

[program:postgres]
command=/webapp/conf/postgres/postgres.sh
directory=/webapp/conf/postgres/
autostart=true
autorestart=true
stopsignal=TERM
priority=300

[program:webapp]
command=/webapp/conf/webapp/webapp.sh
directory=/webapp/conf/webapp/
autostart=true
autorestart=true
stopsignal=TERM
priority=400

[program:nginx]
command=/webapp/conf/nginx/nginx.sh
directory=/webapp/conf/nginx/
autostart=true
autorestart=true
stopsignal=TERM
priority=500
